# templates/fpm/pool.service.erb
[Service]
ExecStart=<%= @fpm_binary %> --fpm-config=/etc/php-fpm.d/<%= @pool %>.conf
User=<%= @user %>
Group=<%= @group %>
Environment="FPM_SOCKETS=/var/run/php-fpm/<%= @pool %>.socket=3"
