# templates/fpm/pool.socket.erb
[Socket]
ListenStream=<%= @listen_port %>

[Install]
WantedBy=sockets.target
