# templates/fpm/pool.socket.erb
[Socket]
<% if @listen_address %>
ListenStream=<%= @listen_address %>:<%= @listen_port %>
<% else %>
ListenStream=<%= @listen_port %>
<% end %>

[Install]
WantedBy=sockets.target
